# 05_ATTACK_T1037.rules
# -----------------------------------------------------------------------
# ATTACK_T1037 - Boot or Logon Initialization Scripts
# -----------------------------------------------------------------------
-w /etc/inittab -p wa -k ATTACK_T1037
-w /etc/init.d/ -p wa -k ATTACK_T1037
-w /etc/init/ -p wa -k ATTACK_T1037
# T1037.004 - RC Scripts
-w /etc/rc.d/ -p wa -k ATTACK_T1037_004
-w /etc/rc.local -p wa -k ATTACK_T1037_004